# oam/install_scripts/mariadb-columnstore.service.in
[Unit]
Description=MariaDB columnstore engine
After=network.target
Wants=mcs-workernode.service

[Service]
Type=oneshot

User=mysql
Group=mysql
LimitNOFILE=65536
LimitNPROC=65536

# This systemd unit is used to manage all MariaDB columnstore systemd units
# Call dummy program
ExecStart=/bin/true

# This service shall be considered active after start
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
