# oam/install_scripts/mcs-loadbrm.service.in
[Unit]
Description=loadbrm
After=mcs-storagemanager.service

[Service]
Type=oneshot

User=mysql
Group=mysql
LimitNOFILE=65536
LimitNProc=65536

ExecStart=/usr/bin/env bash -c "@ENGINE_BINDIR@/mcs-loadbrm.py" 
