# oam/install_scripts/mcs-primproc.service.in
[Unit]
Description=MariaDB columnstore engine PrimProc
PartOf=mcs-workernode.service mcs-controllernode.service
After=mcs-controllernode.service

[Service]
Type=simple

User=mysql
Group=mysql
LimitNOFILE=65536
LimitNProc=65536

ExecStart=/usr/bin/env bash -c "LD_PRELOAD=$(ldconfig -p | grep -m1 libjemalloc | awk '{print $1}') exec @ENGINE_BINDIR@/PrimProc"

Restart=on-failure
TimeoutStopSec=5
