# oam/install_scripts/mcs-controllernode.service.in
[Unit]
Description=MariaDB columnstore engine controllernode
PartOf=mcs-workernode.service
After=mcs-workernode.service

[Service]
Type=simple

User=mysql
Group=mysql
LimitNOFILE=65536
LimitNProc=65536

ExecStart=@ENGINE_BINDIR@/controllernode fg
ExecStop=@ENGINE_BINDIR@/mcs-stop-controllernode.sh $MAINPID
Restart=on-failure
