# oam/install_scripts/mcs-controllernode.service.in
[Unit]
Description=MariaDB columnstore engine controllernode
Wants=primproc.service mcs-writeengineserver.service mcs-exemgr.service
Wants=mcs-dmlproc.service mcs-ddlproc.service

[Service]
Type=simple

User=mysql
Group=mysql
LimitNOFILE=65536
LimitNPROC=65536

ExecStart=@ENGINE_BINDIR@/controllernode fg
ExecStop=@ENGINE_BINDIR@/mcs-stop-controllernode.sh $MAINPID
Restart=on-failure

StandardInput=socket
