# oam/install_scripts/mcs-dmlproc.service.in
[Unit]
Description=MariaDB columnstore engine DMLProc
PartOf=mcs-writeengineserver.service
After=mcs-writeengineserver.service

[Service]
Type=simple

User=mysql
Group=mysql
LimitNOFILE=65536
LimitNProc=65536

ExecStart=@ENGINE_BINDIR@/DMLProc

Restart=on-failure
TimeoutStopSec=5
